# systemd/monitor-garden.service
[Unit]
Description=Garden Monitoring

[Service]
Type=oneshot
ExecStart=/home/pi/software/venv/bin/python /home/pi/software/read.py
Restart=no
