# systemd/monitor-pressure.timer
[Unit]
Description=Pressure Monitoring

[Timer]
AccuracySec=1us
Persistent=true
OnCalendar=*:*:0/10

[Install]
WantedBy=timers.target

