# systemd/monitor-Thermopro-TX2.service
[Unit]
Description=Thermopro-TX2

[Service]
Type=oneshot
ExecStart=/home/pi/software/rtl_433/build/run.sh
Restart=no
